# templates/cloudflared.service.erb
[Unit]
Description=Cloudflared HTTPS DNS proxy
After=network.target

[Service]
User=<%= scope.lookupvar('cloudflared::daemon_user') %>
WorkingDirectory=<%= scope.lookupvar('cloudflared::config_dir') %>
ExecStart=<%= scope.lookupvar('cloudflared::exec_dir') %>/<%= scope.lookupvar('cloudflared::exec_file') %> --config <%= scope.lookupvar('cloudflared::config_file') %> proxy-dns --port 5353
KillMode=process
Restart=on-failure

[Install]
WantedBy=multi-user.target
